# profile.d/ocaml.sh
# eval $(opam config env)
# OPAM configuration
. /Users/martin/.opam/opam-init/init.sh > /dev/null 2> /dev/null || true
