# profile.d/qt.sh
if [[ -d $(brew --prefix qt@5.5) ]]; then
    export PATH="$PATH:$(brew --prefix qt@5.5)/bin"
fi
