# profile.d/roswell.sh
if [[ -d $HOME/.roswell/bin ]]; then
    export PATH=$PATH:$HOME/.roswell/bin
fi
