# profile.d/scala.sh
scala-new() {
    sbt new sbt/scala-seed.g8
}
